# config/systemd/user/niri-exit.service
[Unit]
Description=Run command on Niri exit

[Service]
ExecStart=/usr/bin/pkill --uid "%U" ksecretd
Type=oneshot

[Install]
WantedBy=niri.service
