# systemd-services/tor.service.d/10-periodic.conf
[Service]
Restart=always
RuntimeMaxSec=10m
